# docker/django/services/runserver/finish
#!/usr/bin/with-contenv /bin/sh
exec finish-impostor-victim
