# docker/django/services/uwsgi/finish
#!/usr/bin/with-contenv /bin/sh
exec /usr/local/bin/shutdown-container
